# tests/containerized/image/quadcd-safety.timer
[Unit]
Description=Safety timeout for integration test

[Timer]
OnBootSec=60s

[Install]
WantedBy=timers.target
