# tests/containerized/image/quadcd-safety.service
[Unit]
Description=Safety timeout service

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl poweroff
