# tests/containerized/image/test-sources/hello.container
[Container]
Image=quay.io/podman/hello:latest

[Service]
Restart=no

[Install]
WantedBy=multi-user.target
