# dist/quadcd-sync-user.service
[Unit]
Description=QuadCD Sync Service
After=network-online.target
Wants=network-online.target
ConditionFileIsExecutable=/usr/local/bin/quadcd
ConditionFileNotEmpty=%h/.config/quadcd.toml

[Service]
Type=simple
ExecStart=/usr/local/bin/quadcd sync --service --user
Restart=on-failure
RestartSec=30s

[Install]
WantedBy=default.target
